# quake2-server
#!/bin/sh
# quake2-server
# qui mai  7 16:15:54 EST 1998
#
# chkconfig: 345 98 10
# description: Ativa/Desativa servidor quake2
# Translated to English by fastviper

# Biblioteca de funções
# Functions library
. /etc/rc.d/init.d/functions

GHOME=/usr/share/games/quake2

case "$1" in
  start)
	# Check if the service is already running?
	if [ ! -f /var/lock/subsys/quake2-server ]; then
		msg_starting "Quake2 server"
		busy
		cd $GHOME
		/usr/bin/quake2 +set dedicated 1 +exec server.cfg > /dev/null &
		deltext
		ok
	else
		msg_Already_Running "Quake2 server"
		exit 1
	fi
	touch /var/lock/subsys/quake2-server
	;;
  stop)
	if [ -f /var/lock/subsys/quake2-server ]; then
		msg_stopping "Quake2 server"
		killproc quake2
		rm -f /var/lock/subsys/quake2-server
	else
		msg_Not_Running "Quake2 server"
		exit 1
	fi
	;;
  restart|reload)
	$0 stop
	$0 start
	;;
  status)
	status /usr/bin/quake2
	exit $?
	;;
  *)
	msg_Usage "$0 {start|stop|restart|reload|status}"
	exit 1
esac

exit 0
